# scripts/codeck.service
[Unit]
Description=Codeck - Claude Code Sandbox
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=codeck
Group=codeck
WorkingDirectory=/opt/codeck

Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=production"
Environment="CODECK_DAEMON_PORT=80"
Environment="CODECK_RUNTIME_URL=http://127.0.0.1:7777"
Environment="CODECK_RUNTIME_WS_URL=http://127.0.0.1:7778"
Environment="CODECK_DIR=/workspace/.codeck"
Environment="CODECK_PROJECT_DIR=/opt/codeck"
Environment="CODECK_COMPOSE_FILE=docker/compose.managed.yml"
Environment="HOME=/home/codeck"

ExecStartPre=/usr/bin/docker compose -f /opt/codeck/docker/compose.managed.yml up -d --wait
ExecStart=/usr/bin/node /opt/codeck/apps/daemon/dist/index.js
ExecStopPost=/usr/bin/docker compose -f /opt/codeck/docker/compose.managed.yml down

Restart=always
RestartSec=10
CPUQuota=100%
MemoryMax=512M
NoNewPrivileges=true
ProtectSystem=full
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
